# LaunchAgents/com.tylerlaprade.claude-keepalive.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.tylerlaprade.claude-keepalive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>-lc</string>
		<string>exec "${DOTFILES_DIR:-$HOME/Code/dotfiles}/.claude/idle-keepalive.sh"</string>
	</array>
	<key>StartInterval</key>
	<integer>300</integer>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardOutPath</key>
	<string>/tmp/claude-keepalive.log</string>
	<key>StandardErrorPath</key>
	<string>/tmp/claude-keepalive.log</string>
</dict>
</plist>
